# backend/supervisord.conf
# Combined-image alternative: run the token server AND the agent worker in a
# single container. (docker-compose does NOT use this — it runs the two as
# separate services. This is only for `docker run vma-backend`.)
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:token-server]
command=python token_server.py
directory=/app
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:agent]
command=python agent.py start
directory=/app
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
